# systemd/vpn-policy-routing.service
[Unit]
Description=Apply Codex VPN policy routing
After=network-online.target strongswan-starter.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/vpn-policy-routing.sh
RemainAfterExit=yes
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
